# root/etc/services.d/prosody/run
#!/usr/bin/with-contenv bash

cd /etc/prosody || exit

exec \
	s6-setuidgid abc prosody
